# rhq-411-server/image/start.ini
[supervisord]
nodaemon=true

[program:postgresql]
priority=1
command=runuser -l postgres -c "pg_ctl -D /var/lib/pgsql/data/ -l logfile start"

[program:rhq-storage]
priority=100
command=/rhq-server-4.11.0/bin/rhqctl start --storage

[program:rhq-server]
priority=200
command=/rhq-server-4.11.0/bin/rhqctl start --server

[program:rhq-agent]
priority=200
command=/rhq-server-4.11.0/bin/rhqctl start --agent

[program:sshd]
command=/usr/sbin/sshd -D
